# tools_Vla/CoRunner/com.user.corunner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.corunner</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/uvicorn</string>
    <string>runner_server:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8000</string>
    <string>--workers</string>
    <string>1</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/dunevv/WorkLocal/_AI_/HoudiniElf/tools_Vla/CoRunner</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/corunner.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/corunner.err.log</string>
</dict>
</plist>